# dockerfile-compose/images/supervisor/7.4/supervisord.ini
[supervisord]
nodaemon=true

[program:websocket]
command=php /data/twww/api1/artisan websocket:init
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/websocket.log
stderr_logfile=/server/php/logs/api1/websocket_error.log

[program:223]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=userRisk --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/223.log
stderr_logfile=/server/php/logs/api1/223_error.log

[program:queue_123]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=123 --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/queue_123.log
stderr_logfile=/server/php/logs/api1/queue_123_error.log

[program:queue_1212]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=1212 --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/queue_1212.log
stderr_logfile=/server/php/logs/api1/queue_1212_error.log
