# dockerfile-compose1/images/supervisor/7.4/supervisord.ini
[supervisord]
nodaemon=true

[program:websocket]
command=php /data/twww/api1/artisan websocket:init
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/websocket.log
stderr_logfile=/server/php/logs/api1/websocket_error.log

[program:queue_userRisk]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=userRisk --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/queue_userRisk.log
stderr_logfile=/server/php/logs/api1/queue_userRisk_error.log

[program:queue_giveGiftMoney]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=giveGiftMoney --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/queue_giveGiftMoney.log
stderr_logfile=/server/php/logs/api1/queue_giveGiftMoney_error.log

[program:queue_KsClearBetAmount]
command=php /data/twww/api1/artisan queue:work rabbitmq --queue=KsClearBetAmount --memory=512
autostart=true
autorestart=true
stdout_logfile=/server/php/logs/api1/queue_KsClearBetAmount.log
stderr_logfile=/server/php/logs/api1/queue_KsClearBetAmount_error.log